# vigilantpi.service
[Unit]
Description=Vigilant PI

[Service]
Type=simple
ExecStart=/bin/sh /usr/bin/vigilantpid

[Install]
WantedBy=multi-user.target
