# wifirestarter.service
[Unit]
Description=Wifi Restarter

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/wifirestartd

[Install]
WantedBy=multi-user.target